# src/deb/control/postinst
#!/bin/bash
set -e

cd /usr/share/bayeos-gateway
chmod +x *.sh

# PostgreSQL
service postgresql start 
su - postgres -c /usr/share/bayeos-gateway/pg_postinst.sh

# Gateway 
chmod +x /var/lib/bayeos-gateway/bayeos-gateway.jar
ln -sf /var/lib/bayeos-gateway/bayeos-gateway.jar /etc/init.d/bayeos-gateway
ln -sf /etc/bayeos-gateway/application.properties /var/lib/bayeos-gateway/application.properties

update-rc.d bayeos-gateway defaults
service bayeos-gateway start  

# Apache integration 
a2enmod proxy proxy_http
a2enconf bayeos-gateway.conf
service apache2 restart

